# deployments/systemd/tsbridge.service
[Unit]
Description=Tailscale Bridge Proxy Service
Documentation=https://github.com/jtdowney/tsbridge
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/tsbridge -config /etc/tsbridge/config.toml
Restart=always
RestartSec=10s

# User and group
User=tsbridge
Group=tsbridge

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# State directory management (systemd will create and manage /var/lib/tsbridge)
StateDirectory=tsbridge
StateDirectoryMode=0750

# Allow writing to state directory
ReadWritePaths=/var/lib/tsbridge

# Resource limits
LimitNOFILE=65535
LimitNPROC=512

# Environment
Environment="GOMAXPROCS=0"
EnvironmentFile=-/etc/tsbridge/tsbridge.env

# Logging (systemd will handle this)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tsbridge

# Give it time to start up
TimeoutStartSec=300s
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target